# themes/anonabox/root/etc/uci-defaults/luci-theme-anonabox
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Anonabox=/luci-static/anonabox
	commit luci
EOF
exit 0
